# systemd-nspawn.socket
[Socket]
Accept=true
ListenStream=/binddir/spawnctl

[Install]
WantedBy=sockets.target
